# ubuntu/config/aurora.conf
description "aurora scheduler"

start on started networking
stop on shutdown

setuid aurora

console log

respawn

exec /usr/local/bin/aurora-start
